# .config/systemd/user/theme-apply.path
[Unit]
Description=Watch for wal color changes
After=graphical-session.target


[Path]
PathChanged=%h/.cache/wal/colors.json
Unit=theme-apply.service


[Install]
WantedBy=graphical-session.target