# .config/systemd/user/theme-apply.path
[Unit]
Description=Watch for wal color changes


[Path]
PathChanged=%h/.cache/wal/colors.json


[Install]
WantedBy=default.target