# .config/systemd/user/theme-apply.service
[Unit]
Description=Apply theme when wal colors update


[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/reload-theme